# templates/systemd/backup.service
[Unit]
Description=Backup service placeholder

[Service]
Type=oneshot
ExecStart=/usr/local/bin/backup-all.sh
