# templates/systemd/health-check.timer
[Unit]
Description=Health check timer placeholder

[Timer]
OnBootSec=5m
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
